# templates/example-modules/system-base/rootfs/etc/systemd/system/ssh.service.d/override.conf
[Service]
StandardOutput=journal+console
StandardError=journal+console
ExecStartPre=
ExecStartPre=/bin/sh -c "ls -la /etc/ssh/ /root/.ssh/ /run/sshd && /usr/sbin/sshd -t"
