# initramfs/modules.d/90powernode/persist.mount
[Unit]
Description=Powernode persistent state tmpfs (/persist)
Documentation=https://github.com/nodealchemy/powernode-platform
DefaultDependencies=no
ConditionPathIsMountPoint=!/persist
Before=powernode-agent.service powernode-mount.service local-fs.target

[Mount]
What=tmpfs
Where=/persist
Type=tmpfs
Options=mode=0755,size=64M

[Install]
WantedBy=local-fs.target
